# ansible/roles/wsl/files/mount-z.service
[Unit]
Description=Auto mount Z: drive
After=network.target

[Service]
EnvironmentFile=/etc/default/mount-z
Type=oneshot
ExecStart=/usr/local/bin/mount-z.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
